# ms-weather-bot.service
[Unit]
Description=Mississippi Weather Bot
After=network.target

[Service]
Type=simple
User=weatherbot
WorkingDirectory=/home/weatherbot/ms-weather-bot
ExecStart=/usr/bin/python3 main.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
